# disk_image/create/raspberry_pi_arm64v8/disk_template/rootfs/etc/systemd/system/docker.service.d/options.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/dockerd -H unix:// -H tcp://0.0.0.0:2375
